# templates/nomad.sysv.erb
#!/bin/bash
#
#       /etc/rc.d/init.d/nomad
#
#       Daemonize the nomad agent.
#
# chkconfig:   2345 95 20
# description: Service discovery and configuration made easy. \
#  Distributed, highly available, and datacenter-aware.
# processname: nomad
# pidfile: /var/run/nomad/pidfile

# Source function library.
. /etc/init.d/functions

NOMAD=<%= scope.lookupvar('nomad::bin_dir') %>/nomad
CONFIG=<%= scope.lookupvar('nomad::config_dir') %>
PID_FILE=/var/run/nomad/nomad.pid
LOG_FILE=/var/log/nomad
RPC_ADDR=-rpc-addr=<%= scope.lookupvar('nomad::rpc_addr') %>:<%= scope.lookupvar('nomad::rpc_port') %>

[ -e /etc/sysconfig/nomad ] && . /etc/sysconfig/nomad

export GOMAXPROCS=${GOMAXPROCS:-2}

#
# Create the /var/run/nomad directory, which can live on a tmpfs
# filesystem and be destroyed between reboots.
#
mkrundir() {
        [ ! -d /var/run/nomad ] && mkdir -p /var/run/nomad
        chown <%= scope.lookupvar('nomad::user') %> /var/run/nomad
}

#
# Create a PID file if it doesn't already exist, for clean upgrades
# from previous init-script controlled daemons.
#
KILLPROC_OPT="-p ${PID_FILE}"
mkpidfile() {
        # Create PID file if it didn't exist
        mkrundir
        [ ! -f $PID_FILE ] && pidofproc $NOMAD > $PID_FILE
        chown <%= scope.lookupvar('nomad::user') %> /var/run/nomad
        if [ $? -ne 0 ] ; then
            rm $PID_FILE
            KILLPROC_OPT=""
        fi
}

start() {
        echo -n "Starting nomad: "
        mkrundir
        [ -f $PID_FILE ] && rm $PID_FILE
        daemon --user=<%= scope.lookupvar('nomad::user') %> \
            --pidfile="$PID_FILE" \
            "$NOMAD" agent -pid-file "${PID_FILE}" -config "$CONFIG" <%= scope.lookupvar('nomad::extra_options') %> >> "$LOG_FILE" &
        retcode=$?
        touch /var/lock/subsys/nomad
        return $retcode
}

stop() {
        DELAY=5 # seconds maximum to wait for a leave

        echo -n "Shutting down nomad: "
        mkpidfile

        # If nomad is not acting as a server, exit gracefully
        # Use SIGINT to create a "leave" event, unless the user has explicitly
        # changed that behavior in the Consul config.
        if ("${NOMAD}" info ${RPC_ADDR} 2>/dev/null | grep -q 'server = false' 2>/dev/null) ; then
            nomad_pid=$(cat $PID_FILE)
            killproc $KILLPROC_OPT $NOMAD -INT
            retcode=$?

            # We'll wait if necessary to make sure the leave works, and return
            # early if we can.  If not, escalate to harsher signals.
            try=0
            while [ $try -lt $DELAY ]; do
                if ! checkpid $nomad_pid ; then
                    rm -f /var/lock/subsys/nomad
                    return $retcode
                fi
                sleep 1
                let try+=1
            done
        fi

        # If acting as a server, use a SIGTERM to avoid a leave.
        # This behavior is also configurable.  Avoid doing a "leave" because
        # having servers missing is a bad thing that we want to notice.
        #
        # A SIGTERM will mark the node as "failed" until it rejoins.
        # killproc with no arguments uses TERM, then escalates to KILL.
        killproc $KILLPROC_OPT $NOMAD
        retcode=$?

        rm -f /var/lock/subsys/nomad $PID_FILE
        return $retcode
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        "$NOMAD" info ${RPC_ADDR}
        ;;
    restart)
        stop
        start
        ;;
    reload)
        mkpidfile
        killproc $KILLPROC_OPT $NOMAD -HUP
        ;;
    condrestart)
        [ -f /var/lock/subsys/nomad ] && restart || :
        ;;
    *)
        echo "Usage: nomad {start|stop|status|reload|restart}"
        exit 1
        ;;
esac
retcode=$?
# Don't let the [OK] get stomped on.
echo
exit $retcode
